# scripts/crontab
* * * * * cd /usr/local/var/academic-twitter-csv && ./scripts/job
